# scripts/systemd/forzion-docker-prune.service
[Unit]
Description=Forzion — poda do store Docker (build cache + imagens dangling)
Wants=docker.service
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /opt/forzion/app/scripts/docker-prune.sh
